# .config/systemd/user/drive_backup.service
[Unit]
Description=Backup directory to Google Drive

[Service]
Type=oneshot
ExecStart=%h/dotfiles/extras/rclone/sync-gdrive-verbose.sh
