# .config/systemd/user/drive_backup.timer
[Unit]
Description=Run drive backup every 3 days

[Timer]
OnUnitActiveSec=3d
Persistent=true

[Install]
WantedBy=timers.target
